# dahdi-linux-dkms-git/dahdi.service
[Unit]
Description=DAHDI configuration

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/dahdi_cfg -vvv
ExecStop=/usr/bin/dahdi_cfg -s
